# com.marcelklassen.bauhaus-watcher-local.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.marcelklassen.bauhaus-watcher-local</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/marcelklassen/bauhaus-watcher/venv/bin/python3</string>
        <string>/Users/marcelklassen/bauhaus-watcher/watch_bauhaus.py</string>
        <string>--local-loop</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/marcelklassen/bauhaus-watcher</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/marcelklassen/bauhaus-watcher/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/marcelklassen/bauhaus-watcher/launchd.err.log</string>
</dict>
</plist>
